# EduBlocks/server/edublocks-server.service
[Unit]
Description=EduBlocks Server

[Service]
ExecStart=/opt/edublocks/server/start.sh
Restart=always
User=pi
Group=pi
WorkingDirectory=/opt/edublocks/server

[Install]
WantedBy=multi-user.target
